# etc-remove_me/init.d/S45wpa_supplicant
#!/bin/sh
### BEGIN INIT INFO
# Provides:          wpa_supplicant_wifi
# Required-Start:    $local_fs $network $remote_fs
# Required-Stop:     $local_fs $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start wpa_supplicant for wlan0
# Description:       Start wpa_supplicant for wlan0 at boot
### END INIT INFO

DAEMON=/usr/sbin/wpa_supplicant
NAME=wpa_supplicant_wifi
WLAN_IF=wlan0
CONF=/etc/wpa_supplicant/wpa_supplicant.conf

case "$1" in
  start)
    echo "Starting $NAME"
    ifconfig $WLAN_IF up
    start-stop-daemon --start --quiet \
    --exec $DAEMON -- \
    -B -i $WLAN_IF -c $CONF -Dnl80211,wext;;
  stop)
    echo "Stopping $NAME"
    start-stop-daemon --stop --quiet --exec $DAEMON
    ifconfig $WLAN_IF down
    ;;
  restart|force-reload)
    $0 stop
    sleep 1
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|force-reload}" >&2
    exit 3
    ;;
esac

:
